# docker/etc/rc.local
#!/bin/sh
su-exec app ${APP_HOME}/bin/${APP_NAME} foreground
